# files/root/etc/uci-defaults/luci-tencentddns
#!/bin/sh 

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@ddnspod[-1]
	add ucitrack ddnspod
	set ucitrack.@ddnspod[-1].init=ddnspod
	commit ucitrack
EOF

exit 0
